# scripts/systemd/vx11-madre.service
[Unit]
Description=VX11 Madre
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/elkakas314/vx11
ExecStart=/home/elkakas314/vx11/build/.venv/bin/uvicorn madre.main:app --host 127.0.0.1 --port 52112
Restart=on-failure

[Install]
WantedBy=multi-user.target
[Unit]
Description=VX11 Madre service
After=network.target

[Service]
Type=simple
User=${USER}
WorkingDirectory=/home/elkakas314/vx11
ExecStart=/home/elkakas314/vx11/build/.venv/bin/uvicorn madre.main:app --host 0.0.0.0 --port 52112
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
